# sdca-website-update.cron
# Update repo regularly from Git
0,15,30,45 *  * * *  sdca  cd /var/www/sdca/sdca-website/ && git pull --quiet

